# etc/xdg/xdg-xubuntu/autostart/xcape-super-binding.desktop
[Desktop Entry]
Type=Application
Name=Bind Super Key
Exec=xcape -e 'Super_L=Alt_L|Pause'
OnlyShowIn=XFCE;
