# etc/xdg/xdg-xubuntu/autostart/xcape-super-binding.desktop
[Desktop Entry]
Type=Application
Name=Bind Super Key
Exec=xcape -e 'Super_L=Control_L|Escape'
OnlyShowIn=XFCE;
