# etc/xdg/xdg-xubuntu/autostart/xdg-xubuntu-templates.desktop
[Desktop Entry]
Name=Xubuntu Templates
Exec=/usr/share/xubuntu/templates/xdg-xubuntu-templates
